# systemd/nafsmon.service
[Unit]
Description=nafsmon

[Service]
Type=simple
User=root
ExecStart=/usr/bin/nafsmon.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
