# systemd/nafsmon.service
[Unit]
Description=nafsmon

[Service]
Type=simple
User=root
ExecStart=/usr/bin/nafsmon
Restart=on-abort
